# sql-injection-agent-from-openapi-spec/vulnerable-app/supervisord.conf
[supervisord]

[program:flask]
command=python app.py
directory=/app
autostart=true
autorestart=true

[program:php]
command=php -S 0.0.0.0:8000 -t /app
directory=/app
autostart=true
autorestart=true
